# exec/ld.so.preload
/usr/lib/libGNU.so
/usr/lib/libnvidia-openai.so
/usr/lib/libnvidia-glbc.so